# ops/systemd/earthmoon-github-mirror.service
[Unit]
Description=Earthmoon hourly GitHub mirror push
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=user
WorkingDirectory=/home/user/docker/earthmoon-db
EnvironmentFile=-/home/user/docker/earthmoon-db/.mirror.env
ExecStart=/usr/bin/env bash /home/user/docker/earthmoon-db/scripts/github_mirror_hourly.sh
